# xyz.openbmc_project.NodeInfo.service.in
[Unit]
Description=NodeInfo
After=dbus.service

[Service]
ExecStart=@bindir@/nodeinfo
Restart=always
Type=dbus
BusName=xyz.openbmc_project.NodeInfo

[Install]
WantedBy=multi-user.target